# config/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:memcached]
command=/usr/bin/memcached -u root
